# build/deb/debian/blockbook-zcash.cron.daily
#!/bin/sh
/opt/coins/blockbook/zcash/bin/logrotate.sh
