# build/deb/debian/blockbook-litecoin-testnet.cron.daily
#!/bin/sh
/opt/coins/blockbook/litecoin_testnet/bin/logrotate.sh
