# build/deb/debian/blockbook-dogecoin.cron.daily
#!/bin/sh
/opt/coins/blockbook/dogecoin/bin/logrotate.sh
