# build/deb/debian/blockbook-namecoin.cron.daily
#!/bin/sh
/opt/coins/blockbook/namecoin/bin/logrotate.sh
