# build/deb/debian/blockbook-bch-testnet.cron.daily
#!/bin/sh
find /opt/coins/blockbook/bch-testnet/logs -mtime +6 -type f -delete
