# build/deb/debian/blockbook-denarius.cron.daily
#!/bin/sh
/opt/coins/blockbook/denarius/bin/logrotate.sh
