# build/deb/debian/blockbook-ethereum.cron.daily
#!/bin/sh
/opt/coins/blockbook/ethereum/bin/logrotate.sh
