# build/deb/debian/blockbook-bgold.cron.daily
#!/bin/sh
/opt/coins/blockbook/bgold/bin/logrotate.sh
