# build/deb/debian/blockbook-bcash.cron.daily
#!/bin/sh
/opt/coins/blockbook/bcash/bin/logrotate.sh
