# build/deb/debian/blockbook-groestlcoin.cron.daily
#!/bin/sh
/opt/coins/blockbook/groestlcoin/bin/logrotate.sh
