# build/deb/debian/blockbook-bitcoin.cron.daily
#!/bin/sh
/opt/coins/blockbook/bitcoin/bin/logrotate.sh
