# build/deb/debian/blockbook-btc-testnet.cron.daily
#!/bin/sh
find /opt/coins/blockbook/btc-testnet/logs -mtime +6 -type f -delete
