# build/deb/debian/blockbook-vertcoin.cron.daily
#!/bin/sh
/opt/coins/blockbook/vertcoin/bin/logrotate.sh
