# build/deb/debian/blockbook-digibyte.cron.daily
#!/bin/sh
/opt/coins/blockbook/digibyte/bin/logrotate.sh
