# build/deb/debian/blockbook-btg.cron.daily
#!/bin/sh
/opt/coins/blockbook/btg/bin/logrotate.sh
