# build/deb/debian/blockbook-dash.cron.daily
#!/bin/sh
/opt/coins/blockbook/dash/bin/logrotate.sh
