# build/deb/debian/blockbook-bcash-testnet.cron.daily
#!/bin/sh
/opt/coins/blockbook/bcash_testnet/bin/logrotate.sh
