# build/deb/debian/blockbook-zec.cron.daily
#!/bin/sh
find /opt/coins/blockbook/zec/logs -mtime +6 -type f -delete
