# build/deb/debian/blockbook-bch.cron.daily
#!/bin/sh
find /opt/coins/blockbook/bch/logs -mtime +6 -type f -delete
