# build/deb/debian/blockbook-bitcoin-testnet.cron.daily
#!/bin/sh
/opt/coins/blockbook/bitcoin_testnet/bin/logrotate.sh
