# build/deb/debian/blockbook-bch.cron.daily
#!/bin/sh
/opt/coins/blockbook/bch/bin/logrotate.sh
