# build/deb/debian/blockbook-btc-testnet.cron.daily
#!/bin/sh
/opt/coins/blockbook/btc-testnet/bin/logrotate.sh
