# build/deb/debian/blockbook-zec-testnet.cron.daily
#!/bin/sh
/opt/coins/blockbook/zec-testnet/bin/logrotate.sh
