# build/deb/debian/blockbook-vertcoin-testnet.cron.daily
#!/bin/sh
/opt/coins/blockbook/vertcoin_testnet/bin/logrotate.sh
