# build/deb/debian/blockbook-btc.cron.daily
#!/bin/sh
/opt/coins/blockbook/btc/bin/logrotate.sh
