# build/deb/debian/blockbook-litecoin.cron.daily
#!/bin/sh
/opt/coins/blockbook/litecoin/bin/logrotate.sh
