# build/deb/debian/blockbook-zec.cron.daily
#!/bin/sh
/opt/coins/blockbook/zec/bin/logrotate.sh
