# build/deb/debian/blockbook-ethereum-testnet-ropsten.cron.daily
#!/bin/sh
/opt/coins/blockbook/ethereum_testnet_ropsten/bin/logrotate.sh
