# build/templates/blockbook/debian/cron.daily
{{define "main" -}}
#!/bin/sh
{{.Env.BlockbookInstallPath}}/{{.Coin.Alias}}/bin/logrotate.sh
{{end}}
